# docker/python/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=nobody:nogroup

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[inet_http_server]
port = 9001
username = root
password = root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=root
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d

[program:uwsgi]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=uwsgi --ini app.ini --plugins=python3
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/app/app.log
stopwaitsecs=60
